# roles/opa_setup/templates/opa.container.j2
# OPA Podman Quadlet Container File
# Generated by Ansible - Do not edit manually
# This file is processed by systemd-generator to create a systemd service

[Unit]
Description=Open Policy Agent (OPA) Container
Documentation=https://www.openpolicyagent.org/docs/latest/
After=network-online.target
Wants=network-online.target

[Container]
# Container image
Image={{ opa_container_image }}
Pull={{ opa_pull_policy }}

# Container name
ContainerName={{ opa_container_name }}

{% if opa_container_user %}
# Run as specific user inside container
User={{ opa_container_user }}
{% endif %}

# Volume mounts
{% if opa_use_config_file %}
Volume={{ opa_config_file }}:/config/config.yaml:ro,{{ opa_volume_selinux_label }}
{% endif %}
Volume={{ opa_data_dir }}:/data:{{ opa_volume_selinux_label }}
Volume={{ opa_policy_dir }}:/policies:{{ opa_volume_selinux_label }}
Volume={{ opa_log_dir }}:/logs:{{ opa_volume_selinux_label }}

# Port mappings
{% for port in opa_publish_ports %}
PublishPort={{ port }}
{% endfor %}

# Environment variables
Environment=TZ={{ opa_timezone }}

# OPA command and arguments
{% if opa_use_config_file %}
Exec=run --server --config-file=/config/config.yaml {{ opa_extra_args | join(' ') }}
{% else %}
Exec=run --server --addr={{ opa_server_addr }} {{ opa_extra_args | join(' ') }}
{% endif %}

# Auto-update and restart
AutoUpdate=registry

# Note: Health check removed - OPA container is minimal and doesn't include curl/wget
# Container health can be verified by checking the API endpoint from the host:
#   curl http://localhost:8181/health

[Service]
# Restart policy
Restart={{ opa_restart_policy }}
TimeoutStartSec=60

[Install]
# Start by default
WantedBy=multi-user.target default.target
